# src/etc/init.d/slide-switch.in
#!/bin/sh /etc/rc.common
# Copyright (C) 2016 Jeffery To

prefix=@prefix@
exec_prefix=@exec_prefix@

START=70

boot() {
	"@sbindir@/slide-switch" --init
}
