# src/etc/hotplug.d/button/50-slide-switch.in
#!/bin/sh
# Copyright (C) 2016 Jeffery To

prefix=@prefix@
exec_prefix=@exec_prefix@

"@sbindir@/slide-switch" --button "$BUTTON" &

exit 0
