# etc/hotplug.d/button/50-slide-switch
#!/bin/sh
# Copyright (C) 2016 Jeffery To

/usr/sbin/slide-switch --button "$BUTTON" &

exit 0
